# package/usb-modeswitch/files/gpio.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2007 OpenWrt.org

START=35

start() {
	gpio &
}

stop() {
	killall gpio
}
